# files/crontab
* * * * * echo "docker-phalcon-crontab is running" >> /var/log/cron.log 2>&1
